# Backend/supervisor.conf
[program:my_script]
command=python3 WebServer2.py
directory=/home/runner
autostart=true
autorestart=true
stderr_logfile=/var/log/my_script.err.log
stdout_logfile=/var/log/my_script.out.log
